# root/etc/services.d/aria2b/run
#!/usr/bin/with-contenv bash

if [ "$A2B" == "true" ]
then
  exec aria2b -c "/config/aria2.conf" -u "http://127.0.0.1:${PORT}/jsonrpc" -s $SECRET
fi
