# root/etc/services.d/aria2/run
#!/usr/bin/with-contenv bash

if [[ ! -z $SECRET ]];then
  SECRET_TOKEN="rpc-secret=${SECRET}"
fi

exec \
	s6-setuidgid abc aria2c \
    --conf-path=/config/aria2.conf \
    $SECRET_TOKEN \
    --disk-cache=$CACHE \
    --quiet=$QUIET
  > /dev/stdout \
  2 > /dev/stderr

echo 'Exiting aria2'
