# root/etc/services.d/aria2/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc aria2c \
    --conf-path=/config/aria2.conf \
    --rpc-secret=$SECRET \
    --disk-cache=$Cache
  > /dev/stdout \
  2 > /dev/stderr

echo 'Exiting aria2'